# conf/temp_websocket_server.conf
# supervisor用設定ファイル
[program:temperature_websocket_server]
command=/apps/start_server
numprocs=1
autostart=true
autorestart=true
user=asdf
